# rom-3.215/etc/uci-defaults/50_rtty
. /lib/functions/gl_util.sh

/etc/init.d/rtty disable

uci -q get rtty.general && exit 0

uci set rtty.general=rtty

if [ -f /proc/gl-hw-info/device_ddns ];
then
    id=$(cat /proc/gl-hw-info/device_ddns)
elif [ $(type -q get_default_ddns) ];
then
    id=$(get_default_ddns)
fi

uci set rtty.general.id="$id"

[ -f /etc/wbf/customer ] && {
    uci set rtty.general.ssh_en=1
    uci set rtty.general.web_en=1
}

uci commit rtty
