# rom-3.215/etc/uci-defaults/99_fix_lan_multicast_to_unicast
uci set network.lan.multicast_to_unicast='0'
uci commit
