# rom-3.201/etc/rc.d/S99mount_usdk
#!/bin/sh /etc/rc.common
START=99
mount_usdk()
{
        local devpath
        local action
        if [ -f "/tmp/mount.tmp" ];then
        devpath="$(head -1 /tmp/mount.tmp)"
        action="$(tail -1 /tmp/mount.tmp)"
        env -i DEVPATH=$devpath ACTION=$action /sbin/hotplug-call block
        fi
}

mode_5g_store_cal_data()
{
	. /lib/functions.sh
	. /lib/functions/system.sh
        dd if=/sys/kernel/debug/ieee80211/phy0/ath10k/cal_data of=/lib/firmware/ath10k/cal-pci-0000:00:00.0.bin  bs=2116 count=1
        mac=$(mtd_get_mac_binary art 0)
        mac_5g=$(macaddr_add $mac 2)
	macaddr_2bin $mac_5g | dd of=/lib/firmware/ath10k/cal-pci-0000:00:00.0.bin conv=notrunc bs=1 seek=6 count=6
	echo $mac_5g > /sys/devices/pci0000\:00/0000\:00\:00.0/ieee80211/phy0/macaddress

}

fix_wireless()
{
    . /lib/functions/gl_util.sh
    model=$(get_model)
	[ "$model" = "e750" ] && [ -f "/sys/kernel/debug/ieee80211/phy0/ath10k/cal_data"  ] &&  {
		mode_5g_store_cal_data
	}
	wifi
}

start(){
	fix_wireless
         mount_usdk
        /etc/init.d/mount_usdk disable
}
