# rom-3.201/usr/lib/opkg/info/kmod-shortcut-fe.list
/usr/bin/sfe_dump
/etc/modules.d/09-shortcut-fe
/lib/modules/4.14.221/shortcut-fe.ko
/etc/init.d/shortcut-fe
