# rom-3.216/etc/rc.d/S01init_gpio
#!/bin/sh /etc/rc.common
START=01

set_xe300_usbpower()
{
    . /lib/functions/gl_util.sh
    model=$(get_model)
    #For xe300-no-battery, export gpio16 and turn on the USB power
    [ "$model" = "xe300" ] && [ -d  "/sys/bus/usb/devices/1-1.4/1-1.4:1.0/ttyUSB0/" ] ||  {
		echo 16 >/sys/class/gpio/export
		echo out >/sys/class/gpio/gpio16/direction
		echo 1 >/sys/class/gpio/gpio16/value
    }
}

start(){
	set_xe300_usbpower
}
